# gunicorn/gunicorn.socket
[Unit]
Description=simple django gunicorn socket
[Socket]
ListenStream=/run/gunicorn.sock
[Install]
WantedBy=sockets.target